# backend/deploy/systemd/pruviq-daily-ranking.service
[Unit]
Description=PRUVIQ Daily Strategy Ranking
After=pruviq-api.service
OnFailure=pruviq-alert@%n.service

[Service]
Type=oneshot
User=pruviq
Group=pruviq
WorkingDirectory=/opt/pruviq/current
EnvironmentFile=/opt/pruviq/shared/.env
# RANKING_DIR: `.env` 가 이미 `/opt/pruviq/data/rankings` 로 지정 (archive
# 2026-03~ 보유). systemd `Environment=` 는 `EnvironmentFile=` 을 override
# 하므로 PR #1173 의 shared/ 덮어쓰기가 실제 데이터 경로를 누락시켰음.
# 이제 `.env` 값을 단일 SSoT 로 사용 — 여기서는 주입하지 않음.
# (pruviq-api.service 와 lockstep — 두 프로세스가 같은 .env 를 읽음.)
ExecStart=/opt/pruviq/bin/daily-strategy-ranking.sh
TimeoutStartSec=7200
# 2026-04-18: 3600 → 7200. strategy × direction × timeframe × group 조합이 많고
# 무거운 조합(supertrend/top100)은 개별 60s timeout에 걸려 재시도까지 합치면 1h 초과.
# 병렬화 전까지 제한 완화. 장기적으로는 daily_strategy_ranking.py의 requests를
# async/concurrent.futures 로 전환 검토.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pruviq-daily-ranking
