# backend/deploy/systemd/pruviq-signal-telegram.service
[Unit]
Description=PRUVIQ Telegram Signal Poster (hourly, max 5/hr)
After=network-online.target pruviq-api.service
Wants=network-online.target
OnFailure=pruviq-alert@%n.service

[Service]
Type=oneshot
User=pruviq
Group=pruviq
WorkingDirectory=/opt/pruviq/current
EnvironmentFile=/opt/pruviq/shared/.env
ExecStart=/opt/pruviq/bin/signal-telegram.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120
