# backend/deploy/systemd/pruviq-monitor-full.service
[Unit]
Description=PRUVIQ Health Monitor (hourly full check — site, tunnel, disk)
After=network-online.target pruviq-api.service
Wants=network-online.target
OnFailure=pruviq-alert@%n.service

[Service]
Type=oneshot
User=pruviq
Group=pruviq
WorkingDirectory=/opt/pruviq/current
EnvironmentFile=/opt/pruviq/shared/.env
ExecStart=/opt/pruviq/bin/monitor.sh --full
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60
