# backend/deploy/systemd/pruviq-monitor.timer
[Unit]
Description=PRUVIQ Health Monitor timer (every 5 min)
Requires=pruviq-monitor.service

[Timer]
OnCalendar=*:0/5
RandomizedDelaySec=30
Persistent=true
Unit=pruviq-monitor.service

[Install]
WantedBy=timers.target
