# backend/deploy/systemd/pruviq-signal-telegram.timer
[Unit]
Description=PRUVIQ Telegram Signal timer (every hour at :05)
Requires=pruviq-signal-telegram.service

[Timer]
OnCalendar=*:05
RandomizedDelaySec=30
Persistent=true
Unit=pruviq-signal-telegram.service

[Install]
WantedBy=timers.target
