# backend/deploy/systemd/pruviq-daily-ranking.timer
[Unit]
Description=PRUVIQ Daily Strategy Ranking timer (daily 20:00 UTC = 05:00 KST+1)
Requires=pruviq-daily-ranking.service

[Timer]
OnCalendar=*-*-* 20:00:00
RandomizedDelaySec=30
Persistent=true
Unit=pruviq-daily-ranking.service

[Install]
WantedBy=timers.target
