# backend/deploy/systemd/pruviq-commit-data.service
[Unit]
Description=PRUVIQ Data Commit + Push (triggered after refresh-data)
# No timer of its own — it's called via ExecStartPost from refresh-data,
# so it only runs when there's something to commit.
OnFailure=pruviq-alert@%n.service

[Service]
Type=oneshot
User=pruviq
Group=pruviq
WorkingDirectory=/opt/pruviq/current
EnvironmentFile=/opt/pruviq/shared/.env
# Make the deploy key visible
Environment="GIT_SSH_COMMAND=ssh -i /opt/pruviq/.ssh/id_ed25519 -o IdentitiesOnly=yes"
ExecStart=/opt/pruviq/bin/commit-data.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60
