# backend/deploy/systemd/pruviq-update-ohlcv.timer
[Unit]
Description=PRUVIQ OHLCV Update timer (hourly at :15)
Requires=pruviq-update-ohlcv.service

[Timer]
# Fires every hour at :15 past. Changed from "every 4h" on 2026-04-21.
#
# Why hourly is the practical maximum (not faster):
#   PRUVIQ backtests use 1-hour candles (*_1h.csv). A new closed bar is
#   only published by OKX once per hour on the :00 boundary. Refreshing
#   at 30m/15m/5m intervals receives the same candle multiple times —
#   pure API waste, zero data-quality gain.
#
# Why :15 past (not :05):
#   OKX typically publishes the just-closed bar within 1-3 minutes, but
#   occasionally lags up to 5. Firing at :15 gives a 15-min safety margin
#   so we never chase an unpublished bar.
#
# Rate-limit envelope (well under):
#   market/candles endpoint: 40 req/2s per IP = 20 req/s sustained.
#   240 symbols × 1 req = ~12s of API time per run.
#   Hourly = 24 runs/day × ~30s = 12 min/day of API activity. Trivial.
OnCalendar=*-*-* *:15:00
RandomizedDelaySec=60
Persistent=true
Unit=pruviq-update-ohlcv.service

[Install]
WantedBy=timers.target
