# backend/deploy/systemd/pruviq-monitor-full.timer
[Unit]
Description=PRUVIQ Health Monitor Full timer (hourly)
Requires=pruviq-monitor-full.service

[Timer]
OnCalendar=hourly
RandomizedDelaySec=30
Persistent=true
Unit=pruviq-monitor-full.service

[Install]
WantedBy=timers.target
