# backend/deploy/systemd/pruviq-daily-ranking.timer
[Unit]
Description=PRUVIQ Daily Strategy Ranking timer (daily 01:00 UTC = 10:00 KST)
Requires=pruviq-daily-ranking.service

[Timer]
OnCalendar=*-*-* 01:00:00
RandomizedDelaySec=30
Persistent=true
Unit=pruviq-daily-ranking.service

[Install]
WantedBy=timers.target
